# init/opticon-api.service
[Unit]
Description=Opticon API

[Service]
Type=oneshot
ExecStart=/usr/sbin/opticon-api

[Install]
WantedBy=multi-user.target
