# init/opticon-collector.service
[Unit]
Description=Opticon Collector

[Service]
Type=oneshot
ExecStart=/usr/sbin/opticon-collector

[Install]
WantedBy=multi-user.target
